# web/src/bin/findora-web.service
[Unit]
Description=Findora web daemon
After=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=1
User=servicefindora
WorkingDirectory=/home/servicefindora/findora-toolbox/web
ExecStart=python3 /home/servicefindora/findora-toolbox/web/src/app.py
SyslogIdentifier=findora-toolbox-web
StartLimitInterval=0
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target